# source/etc/rc.local
#!/bin/sh -e
#
# rc.local
#

cd /opt/word-changer
source .venv/bin/activate
./word-changer-web > /var/log/word-changer-web.log 2>/var/log/word-changer-web.errors &
word-changer > /var/log/word-changer.log 2>/var/log/word-changer.errors &

exit 0
